# systemd/inaki.service
[Unit]
Description=Iñaki — asistente personal agentico
Documentation=https://github.com/inaki
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/inaki
ExecStart=/home/pi/inaki/.venv/bin/python main.py --daemon
ExecReload=/bin/kill -HUP $MAINPID

# Reiniciar automáticamente si falla (salvo salida limpia o señal de stop)
Restart=on-failure
RestartSec=5s

# Logs van a journald (ver con: journalctl -u inaki -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inaki

# Límites de recursos para Raspberry Pi 5 (4GB RAM)
# Ajustar según uso real
MemoryMax=2G
TasksMax=64

# Variables de entorno (alternativa a config/global.secrets.yaml)
# EnvironmentFile=/home/pi/inaki/.env

# Shutdown gracioso: SIGTERM, esperar 30s, luego SIGKILL
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
